# third_party/tuntap/tuntap-20141104/tuntap/pkg/launchd/net.sf.tuntaposx.tap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>net.sf.tuntaposx.tun</string>
    <key>ProgramArguments</key>
    <array>
        <string>/sbin/kextload</string>
        <string>/Library/Extensions/tun.kext</string>
    </array>
    <key>KeepAlive</key>
    <false/>
    <key>RunAtLoad</key>
    <true/>
    <key>UserName</key>
    <string>root</string>
</dict>
</plist>
